# systemd/claude-codex-sync-aws.service
[Unit]
Description=Back up Claude Code and Codex CLI chat history to AWS S3
Documentation=https://github.com/yunwei37/claude-codex-sync-aws

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/claude-codex-sync-aws/env
ExecStart=%h/.local/bin/claude-codex-sync-aws backup
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
