# systemd/claude-codex-sync-aws.timer
[Unit]
Description=Weekly Claude/Codex chat history backup

[Timer]
OnCalendar=Sun *-*-* 03:30:00
Persistent=true
RandomizedDelaySec=2h
Unit=claude-codex-sync-aws.service

[Install]
WantedBy=timers.target
